# deployments/smartmeter-exporter.service
[Unit]
Description=Smartmeter exporter

[Service]
ExecStart=/opt/u-one/echonetlite/start.sh

[Install]
WantedBy=multi-user.target
